# infrastructure/scripts/axontis-horizon.conf
; ============================================
; Configuration Supervisor - Laravel Horizon
; ============================================
;
; Horizon gère TOUS les workers automatiquement.
; Un seul processus Supervisor suffit — Horizon s'occupe du reste :
;   - Auto-scaling des workers par queue
;   - Monitoring, métriques, retry
;   - Dashboard web accessible sur /horizon
;
; MIGRATION depuis les anciens workers :
;   1. sudo supervisorctl stop axontis-workers:*
;   2. sudo mv /etc/supervisor/conf.d/axontis-worker.conf /etc/supervisor/conf.d/axontis-worker.conf.disabled
;   3. sudo mv /etc/supervisor/conf.d/axontis-workers.conf /etc/supervisor/conf.d/axontis-workers.conf.disabled
;   4. sudo cp axontis-horizon.conf /etc/supervisor/conf.d/
;   5. sudo supervisorctl reread
;   6. sudo supervisorctl update
;   7. sudo supervisorctl start axontis-horizon
;
; Commandes utiles :
;   sudo supervisorctl status
;   sudo supervisorctl restart axontis-horizon
;   sudo supervisorctl tail -f axontis-horizon
;
; ============================================

[program:axontis-horizon]
process_name=%(program_name)s
command=php /var/www/axontis/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/axontis/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

