# infrastructure/scripts/axontis-reverb.conf
# ─────────────────────────────────────────────────────────
# Laravel Reverb — WebSocket server pour broadcasting alarme
# ─────────────────────────────────────────────────────────
# Ce processus sert les connexions WebSocket pour les
# événements temps réel (alarme, arm/disarm, connection).
# Canal : private-installation.{uuid}
# ─────────────────────────────────────────────────────────

[program:axontis-reverb]
command=php /var/www/axontis/artisan reverb:start --port=8080
directory=/var/www/axontis
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/axontis/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=30

