# infrastructure/scripts/supervisor-workers.conf
; ============================================
; Configuration Supervisor pour les Workers de Notifications
; ============================================
;
; Ce fichier configure les workers de queue Laravel pour gérer
; les notifications par canal de communication.
;
; Installation :
;   1. Copiez ce fichier dans /etc/supervisor/conf.d/
;   2. Ajustez les chemins selon votre installation
;   3. Exécutez : sudo supervisorctl reread
;   4. Exécutez : sudo supervisorctl update
;   5. Exécutez : sudo supervisorctl start all
;
; Commandes utiles :
;   sudo supervisorctl status
;   sudo supervisorctl restart all
;   sudo supervisorctl tail -f axontis-worker-emails
;
; ============================================

; --------------------------------------------
; Worker pour les Emails (Priorité Haute)
; --------------------------------------------
; Gère : WelcomeNotification, InvoiceNotification, etc.
; Failover automatique : Resend → Mailgun → Brevo
;
[program:axontis-worker-emails]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/axontis/artisan queue:work redis --queue=emails --tries=3 --timeout=60 --sleep=3 --max-jobs=1000 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/axontis/storage/logs/worker-emails.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; --------------------------------------------
; Worker pour les SMS
; --------------------------------------------
; Gère : SmsVerificationNotification, alertes SMS
; Provider : Twilio / Vonage
;
[program:axontis-worker-sms]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/axontis/artisan queue:work redis --queue=sms --tries=3 --timeout=30 --sleep=3 --max-jobs=500 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/axontis/storage/logs/worker-sms.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=60

; --------------------------------------------
; Worker pour WhatsApp
; --------------------------------------------
; Gère : WhatsAppOrderUpdateNotification
; Provider : Twilio WhatsApp Business
;
[program:axontis-worker-whatsapp]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/axontis/artisan queue:work redis --queue=whatsapp --tries=3 --timeout=30 --sleep=3 --max-jobs=500 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/axontis/storage/logs/worker-whatsapp.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=60

; --------------------------------------------
; Worker pour Telegram et Slack
; --------------------------------------------
; Gère : UrgentAlertNotification, notifications Telegram
; Canaux : Telegram Bot, Slack Webhooks
;
[program:axontis-worker-messaging]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/axontis/artisan queue:work redis --queue=telegram --tries=3 --timeout=30 --sleep=3 --max-jobs=500 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/axontis/storage/logs/worker-messaging.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=60

; --------------------------------------------
; Worker par défaut (autres jobs)
; --------------------------------------------
; Gère : tous les autres jobs non spécifiques
;
[program:axontis-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/axontis/artisan queue:work redis --queue=default --tries=3 --timeout=90 --sleep=3 --max-jobs=1000 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/axontis/storage/logs/worker-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=90

; --------------------------------------------
; Groupe de tous les workers
; --------------------------------------------
[group:axontis-workers]
programs=axontis-worker-emails,axontis-worker-sms,axontis-worker-whatsapp,axontis-worker-messaging,axontis-worker-default
priority=999
